# target/linux/qualcommax/ipq60xx/base-files/etc/uci-defaults/01_promisc_mode
#!/bin/sh

. /lib/functions.sh

board=$(board_name)

case "$board" in
yuncore,fap650)
	uci set network.wan_dev=device
	uci set network.wan_dev.name='wan'
	uci set network.wan_dev.promisc='1'
	uci commit
	reload_config
esac
